# packages/gaming-station/host/systemd-user/sunshine-host.service
[Unit]
Description=Host Xorg + Sunshine Session

[Service]
Type=simple
Environment=SUNSHINE_BIN={{SUNSHINE_BIN}}
ExecStart={{HOST_START_SCRIPT_PATH}}
Restart=on-failure
RestartSec=2
KillMode=mixed

[Install]
WantedBy=default.target
